# ghost-analytics.service
[Unit]
Description=Ghost Analytics - Hourly Report Generator
After=poly-latency-bot.service

[Service]
Type=oneshot
WorkingDirectory=/root/.gemini/antigravity/scratch/poly
ExecStart=/root/.gemini/antigravity/scratch/poly/poly_venv/bin/python /root/.gemini/antigravity/scratch/poly/analytics_engine.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
